# examples/docker/aurora-run/supervisord.conf
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#
[supervisord]
nodaemon=true
childlogdir=/var/log/aurora

[program:zookeeper]
environment=JVMFLAGS="-Djava.net.preferIPv4Stack=true"
command=/usr/share/zookeeper/bin/zkServer.sh start-foreground

[program:mesos-master]
environment=LD_LIBRARY_PATH=/usr/lib/jvm/java-8-openjdk-amd64/jre/lib/amd64/server
command=/usr/sbin/mesos-master
    --zk=zk://127.0.01:2181/mesos
    --work_dir=/usr/local/aurora/master/db
    --quorum=1

[program:mesos-slave]
environment=LD_LIBRARY_PATH=/usr/lib/jvm/java-8-openjdk-amd64/jre/lib/amd64/server
command=bash -c 'sleep 5 && /usr/sbin/mesos-slave \
  --master=zk://127.0.0.1:2181/mesos \
  --hostname=127.0.0.1 \
  --attributes="host:127.0.0.1;rack:a" \
  --work_dir="/var/lib/mesos" \
  --containerizers=docker,mesos'

# TODO(dnorris) add authentication support
[program:aurora-scheduler]
environment=GLOG_v="0",LIBPROCESS_PORT="8083",DIST_DIR="/aurora/dist"

command=bash -c 'JAVA_OPTS=\
  "-Djava.library.path=/usr/lib -Dlog4j.configuration=\"file:///etc/zookeeper/conf/log4j.properties\"" \
  /usr/bin/aurora-scheduler \
    -cluster_name=devcluster \
    -http_port=8081 \
    -native_log_quorum_size=1 \
    -zk_endpoints=localhost:2181 \
    -mesos_master_address=zk://localhost:2181/mesos \
    -serverset_path=/aurora/scheduler \
    -native_log_zk_group_path=/aurora/replicated-log \
    -native_log_file_path=/var/db/aurora \
    -backup_dir=/var/lib/aurora/backups \
    -thermos_executor_path=$DIST_DIR/thermos_executor.pex \
    -thermos_executor_flags="--announcer-enable --announcer-ensemble localhost:2181" \
    -vlog=INFO \
    -logtostderr \
    -allowed_container_types=DOCKER,MESOS \
    -use_beta_db_task_store=true'

[program:thermos_observer]
command=/aurora/dist/thermos_observer.pex
  --root=/var/run/thermos
  --port=1388
  --log_to_disk=NONE
  --log_to_stderr=google:INFO
